# stemcell_builder/stages/bosh_agent/assets/runit/agent/run
#!/bin/bash
#
# Copyright (c) 2009-2012 VMware, Inc.

set -e

export PATH=/var/vcap/bosh/bin:$PATH
exec 2>&1

# /usr/bin/nice -n -10
# /bin/nice on centos
exec /var/vcap/bosh/bin/bosh_agent -c -P centos -I $(cat /etc/infrastructure)
